# nza-userland/components/clamav/clamav-daemon.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='clamav'>

<service
	name='network/clamav'
	type='service'
	version='1'>

	<create_default_instance enabled='true' />

	<dependency
	    name='fs-local'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

	<dependency
	    name='network-service'
	    grouping='require_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/network/service' />
	</dependency>

	<dependency
	    name='name-services'
	    grouping='require_all'
	    restart_on='refresh'
	    type='service'>
		<service_fmri value='svc:/milestone/name-services' />
	</dependency>

	<dependency
	    name='identity'
	    grouping='optional_all'
	    restart_on='refresh'
	    type='service'>
		<service_fmri value='svc:/system/identity:domain' />
	</dependency>

	<dependency
	    name='system-log'
	    grouping='optional_all'
	    restart_on='none'
	    type='service'>
		<service_fmri value='svc:/system/system-log' />
	</dependency>

	<dependency
	    name='config-file'
	    grouping='require_all'
	    restart_on='refresh'
	    type='path'>
		<service_fmri
		    value='file://localhost/etc/clamav/clamd.conf' />
		<service_fmri
		    value='file://localhost/etc/clamav/freshclam.conf' />
	</dependency>

	<dependent
		name='clamav_multi-user'
		grouping='optional_all'
		restart_on='none'>
		<service_fmri value='svc:/milestone/multi-user' />
	</dependent>

	<exec_method
		type='method'
		name='start'
		exec='/lib/svc/method/clamav start'
		timeout_seconds='60' />

	<exec_method
		type='method'
		name='stop'
		exec='/lib/svc/method/clamav stop %{restarter/contract}'
		timeout_seconds='60' />

	<exec_method
		type='method'
		name='refresh'
		exec=':kill -HUP'
		timeout_seconds='30' />

	<property_group name='general' type='framework'>
		<propval name='action_authorization' type='astring'
			value='solaris.smf.manage.clamav' />
	</property_group>

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
				ClamAV Anti-Virus
			</loctext>
		</common_name>
	</template>

</service>

</service_bundle>
