# nza-userland/components/c-icap/c-icap.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<service_bundle type='manifest' name='export'>
  <service name='network/c-icap' type='service' version='0'>
    <single_instance/>
    <dependency name='fs' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local'/>
    </dependency>
    <dependency name='net' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/network/loopback'/>
    </dependency>
    <exec_method name='start' type='method' exec='/lib/svc/method/svc-c-icap start' timeout_seconds='18446744073709551615'>
      <method_context/>
    </exec_method>
    <exec_method name='stop' type='method' exec='/lib/svc/method/svc-c-icap stop' timeout_seconds='18446744073709551615'>
      <method_context/>
    </exec_method>
    <exec_method name='restart' type='method' exec='/lib/svc/method/svc-c-icap restart' timeout_seconds='18446744073709551615'>
      <method_context/>
    </exec_method>
    <instance name='default' enabled='true'>
      <property_group name='general' type='framework'>
          <propval name='action_authorization' type='astring'
              value='solaris.smf.manage.cicap' />
        <property name='package' type='astring'/>
      </property_group>
    </instance>
  </service>
</service_bundle>

