# src/ncp/nexenta-ldconfig/debian/sysv-rc.install
usr/sbin/update-rc.d
usr/sbin/invoke-rc.d
etc/init.d/rcS
etc/default/rcS
